# ansible/files/ebs-playground.service
[Unit]
Description=The EBS Playground service

[Service]
ExecStart=/home/ubuntu/ebs-playground

[Install]
WantedBy=default.target